# templates/borg-backup.timer.epp
<%- | Hash[String,String] $backuptime,
| -%>
# THIS FILE IS MANAGED BY PUPPET
[Unit]
Description=Daily Timer(s) for borg-backup

[Timer]
# OnBootSec=10min
# OnUnitActiveSec=1d
<% include stdlib -%>
<% $backuptime.keys.each | $key | { -%>
# Daily named job "<%= $key %>" at <%= $backuptime[$key] %>
OnCalendar=*-*-* <%= $backuptime[$key] %>
<% } -%>
Persistent=True

[Install]
WantedBy=multi-user.target
