# templates/borg-backup.timer.epp
# THIS FILE IS MANAGED BY PUPPET
[Unit]
Description=Daily Timer for borg-backup

[Timer]
# OnBootSec=10min
# OnUnitActiveSec=1d
# Daily at 8am
OnCalendar=*-*-* 18:30:00
Persistent=True

[Install]
WantedBy=multi-user.target
